# deployment/phone-services.service
[Unit]
Description=Yealink Phone Services Django application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yealink
Group=yealink
WorkingDirectory=/opt/yealinkService
EnvironmentFile=/etc/yealinkService/phone-services.env
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/yealinkService/diversion.log
StandardError=append:/var/log/yealinkService/diversion.log
ExecStart=/opt/yealinkService/.venv/bin/gunicorn --bind 127.0.0.1:8001 --workers 2 --timeout 60 --access-logfile /var/log/yealinkService/gunicorn-access.log --capture-output yealinkService.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
